# unbound/sv/unbound/finish
#!/bin/sh

pidfile="/usr/local/etc/unbound/unbound.pid"
echo "[unbound] Stopping unbound..."

if [ -f $pidfile ]; then
	rm $pidfile
fi
