# dnsdist/init.d/dnsdist
#!/sbin/openrc-run
################################################################################
## init.d/dnsdist
################################################################################
name="dnsdist"
supervisor="supervise-daemon"

extra_commands="checkconfig"

description="A high-performance, programmable, secure, load-balancing DNS proxy"

: ${cfgfile:=${DNSDIST_CONFIG:-"/usr/local/etc/dnsdist.conf"}}

command="/usr/local/bin/dnsdist"
command_args="-u dnsdist -g dnsdist -C $cfgfile"
command_args_foreground="--supervised"

required_files="$cfgfile"

depend() {
	need net dcron rootwyrm
	use logger
	after firewall 
	before opennhrp bgpd
	provide dns
}

checkconfig() {
	ebegin "Checking $name config file"
	$command --check-config -C $cfgfile > /dev/null
	eend $?
}

start_pre() {
	/opt/rootwyrm/bin/identify.sh
	if [ -f /message ]; then
		cat /message
	fi
	checkconfig
}

start()
{
	ebegin "Starting dnsdist"
	exec $command $command_args > /dev/null &
	eend $?
}
