# dnsdist/init.d/dnsdist
#!/sbin/openrc-run
################################################################################
## init.d/dnsdist
################################################################################
name="dnsdist"

extra_commands="checkconfig"

description="A high-performance, programmable, secure, load-balancing DNS proxy"

#: ${cfgfile:=${DNSDIST_CONFIG:-"/usr/local/etc/dnsdist/dnsdist.conf"}}

command="/usr/local/bin/dnsdist"
command_args="-u dnsdist -g dnsdist -C $cfgfile --supervised"

required_files="$cfgfile"

depend() {
	need net dcron
	use logger
	after firewall 
	before opennhrp bgpd
	provide dns
}

checkconfig() {
	ebegin "Checking $name config file"
	$command --check-config -C $cfgfile > /dev/null
	eend $?
}

start_pre() {
	/opt/rootwyrm/bin/firstboot
	if [ $? -ne 0 ]; then
		exit $?
	fi
	: ${cfgfile:=${DNSDIST_CONFIG:-"/usr/local/etc/dnsdist/dnsdist.conf"}}
	checkconfig
}

start()
{
	ebegin "Starting dnsdist"
	exec $command $command_args > /dev/null &
	eend $?
}
