# base/service/cron/finish
#!/bin/sh

kill -SIGTERM $(cat /var/run/dcron.pid)

if [ -f /var/run/dcron.pid ]; then
	rm /var/run/dcrond.pid
fi
