# container/nsd/init.d/nsd
#!/sbin/openrc-run
################################################################################
## init.d/nsd
################################################################################
supervisor=supervise-daemon

extra_commands="checkconfig"
extra_started_commands="reload reopen"

name="NSD"

description="An authoritative-only, high performance, open source name server"
description_checkconfig="Checks NSD configuration file for errors"
description_reload="Reloads NSD config file and modified zone files from disk"
description_reopen="Reopen NSD log file (for logrotate)"

: ${cfgfile:=${NSD_CONFIG:-"/usr/local/etc/nsd/nsd.conf"}}
: ${user:=${NSD_USER:-"dnsdocker"}}

command="/usr/local/sbin/nsd"
command_args="-c $cfgfile -l /var/log/nsd.log -u ${user} ${command_args:-}"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
start_stop_daemon_args="--wait 50"

required_files="$cfgfile"

depend() {
	need net dcron rootwyrm
	use logger
}

checkconfig() {
	ebegin "Checking $name config file"
	/usr/local/sbin/nsd-checkconf "$cfgfile"
	eend $?
}

reload() {
	ebegin "Reloading $name config file and all zones"
	/usr/local/sbin/nsd-control reconfig && /usr/local/sbin/nsd-control reload
	eend $?
}

reopen() {
	ebegin "Reopening $name log file"
	/usr/local/sbin/nsd-control log_reopen
	eend $?
}

start_pre() {
    checkconfig
	/opt/rootwyrm/bin/identify.sh
	if [ -f /message ]; then
		cat /message
	fi
}

start()
{
	ebegin "Starting NSD"
	$command $command_args
	eend $?
}
	
stop() {
	ebegin "Writing out NSD zones" 
	/usr/local/sbin/nsd-control write > /dev/null
	eend $?
	ebegin "Stopping NSD"
	/usr/local/sbin/nsd-control stop > /dev/null
	eend $?
}

# vim ft=sh:sw=4:ts=4:expandtab
