# nsd/sv/nsd/finish
#!/bin/sh
################################################################################
## sv/nsd/finish
################################################################################

echo "Stopping nsd..."
