# base/service/firstboot/finish
#!/bin/sh
mv /etc/service/firstboot /etc/sv
